# doc/pets.init
#! /bin/sh
#
### BEGIN INIT INFO
# Provides:          pets
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: P.E.T.S.
# Description:       persistent encryptable tmpfs server
### END INIT INFO

